# pulseaudio.socket
[Unit]
Description=Sound System (System)

[Socket]
ListenStream=%t/pulse/native

[Install]
WantedBy=sockets.target
